# uwsgi_autostart.sh
#!/bin/bash
# chkconfig: 2345 10 90
# description: myservice ....
/usr/bin/uwsgi --ini /opt/api-server/uwsgi.ini